# rc.local
#!/bin/bash -e
perl ~/perlircssl.pl &
exit 0
